# debian/DEBIAN/postinst
#!/bin/bash

#
#		Postinstall script
#

# Enable service
systemctl enable bootsetup
systemctl enable bootsetup-docker