# debian/DEBIAN/prerm
#!/bin/bash

#
#		Pre-removal script
#
systemctl disable bootsetup
systemctl disable bootsetup-docker